# nightscout/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Nightscout
# Runs the Nginx daemon
# ==============================================================================

# Wait for Nightscout to become available
bashio::log.info "Nginx waiting for Nightscout..."
bashio::net.wait_for 1337

bashio::log.info "Starting NGINX..."
if bashio::config.true 'ssl'; then
    options='/etc/nginx/nginx-ssl.conf'
else
    options='/etc/nginx/nginx.conf'
 fi

exec nginx -c "$options"
